# supervisor/supervisord.conf

[program:controller]
command=/usr/local/bin/gunicorn -c /code/backup-mgmt/gunicorn_conf.py --chdir /code/backup-mgmt/backup backup.wsgi:application --reload
