# supervisor/supervisord.conf
[program:controller]
command=/usr/local/bin/gunicorn -c /home/trangnth/backup-mgmt/gunicorn_conf.py --chdir /home/trangnth/backup-mgmt/backup backup.wsgi:application --reload
