# supervisor/supervisord.conf
[program:controller]
command=/usr/local/bin/gunicorn -c /code/backup_mgmt/gunicorn_conf.py --chdir /code/backup_mgmt/backup backup.wsgi:application --reload
